# ansible/roles/pivx_instance/templates/pivxd@.service.j2
# /etc/systemd/system/pivxd@{{ instance.name }}.service
# Generated by Ansible — do not edit manually.
# systemd instance unit for PIVX instance: {{ instance.name }}
#
# Binary: /opt/pivx/current/bin/pivxd  (→ /opt/pivx/releases/{{ pivx_version }}/)
# Protocol: {{ instance.protocol_class }} | Role: {{ instance.role }} | Phase: {{ lifecycle_phase }}
#
# Useful commands:
#   systemctl status  pivxd@{{ instance.name }}
#   journalctl -u     pivxd@{{ instance.name }} -f
#   pivx-cli -conf={{ instance.confdir }}/pivx.conf getblockchaininfo

[Unit]
Description=PIVX Daemon — {{ instance.name }} ({{ instance.protocol_class }}/{{ instance.role }})
Documentation=https://github.com/PIVX-Project/PIVX
After=network-online.target{% if instance.protocol_class == 'tor' %} tor.service{% endif %}

Wants=network-online.target
{% if instance.protocol_class == 'tor' %}
Requires=tor.service
{% endif %}

[Service]
Type=forking
EnvironmentFile={{ instance.confdir }}/instance.env
User={{ pivx_user }}
Group={{ pivx_group }}

ExecStart={{ pivx_current_link }}/bin/pivxd \
  -conf={{ instance.confdir }}/pivx.conf \
  -datadir={{ instance.datadir }} \
  -daemon \
  -pid={{ instance.datadir }}/pivxd.pid

ExecStop={{ pivx_current_link }}/bin/pivx-cli \
  -conf={{ instance.confdir }}/pivx.conf \
  -datadir={{ instance.datadir }} \
  stop

PIDFile={{ instance.datadir }}/pivxd.pid

# Restart policy
Restart=on-failure
RestartSec=30s
TimeoutStartSec=120s
TimeoutStopSec=90s

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths={{ instance.datadir }} {{ instance.logdir }}

# Logging goes to journald; Vector picks it up from there
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pivxd-{{ instance.name }}

[Install]
WantedBy=multi-user.target

